# septic.service.in
[Unit]
Description=Septic: the Protab contest daemon
After=network.target

[Service]
Type=simple
ExecStart=${install_dir}/septic -s
Restart=always
User=septic

[Install]
WantedBy=multi-user.target
